# docker/supervisord.conf
[supervisord]
nodaemon=true

[fcgi-program:uvicorn]
socket=tcp://localhost:8000
command=uvicorn --fd 0 base:app
numprocs=4
process_name=uvicorn-%(process_num)d
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:webserver]
command=/usr/sbin/nginx -g "daemon off;" -c /etc/nginx/nginx.conf
priority=900
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
username=www-data
autorestart=true
